# recipes-core/initscripts/initscripts/rc.local
#!/bin/sh
modprobe spi_bcm2835
modprobe fbtft_device name=ft81x
