# recipes-core/hwclock-load/hwclock-load/rtc-exists.path
[Unit] 
Description=Wait for /dev/rtc to appear
PartOf=hwclock-load.service

[Path]
PathExists=/dev/rtc

[Install]
WantedBy=multi-user.target

