# recipes-core/frame-buffer/fb1-1.1/fb1.path
[Unit]
Description=Wait for /dev/fb1 to exist

[Path]
PathExists=/dev/fb1
# Unit=fb1.service

[Install]
WantedBy=sysinit.target

